# configs/systemd/mongo_backup.service
[Unit]
Description=MongoDB Backup Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mongo_backup.sh
User=root
