# configs/systemd/recolecta_info.service
[Unit]
Description=Recolecta info script

[Service]
Type=oneshot
ExecStart=/usr/local/bin/recolecta_info.sh
User=ubuntu
