# configs/systemd/recolecta_info.timer
[Unit]
Description=Run recolecta_info every 5 minutes

[Timer]
OnCalendar=*:0/5
Unit=recolecta_info.service

[Install]
WantedBy=timers.target
